# packaging/systemd/tapx-panel.service
[Unit]
Description=TapX control panel and local runtime manager
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/tapx/tapx.env
ExecStart=/usr/local/bin/tapx-panel -listen=${TAPX_PANEL_LISTEN} -base-path=${TAPX_PANEL_BASE_PATH}
Restart=on-failure
RestartSec=2s
LimitNOFILE=1048576

# The first Linux service runs as root because TapX creates TAP/TUN devices,
# configures addresses/routes/iptables, and starts raw sockets. A later
# hardening pass can split privileges after the netapply boundary stabilizes.
User=root
Group=root
RuntimeDirectory=tapx
StateDirectory=tapx
LogsDirectory=tapx
ReadWritePaths=/etc/tapx /run/tapx /var/lib/tapx /var/log/tapx
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=false
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SETPCAP CAP_DAC_OVERRIDE CAP_CHOWN CAP_FOWNER

[Install]
WantedBy=multi-user.target
